# launchd/com.hashiru.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hashiru.morning</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/hachikenkazuki/hashiru/hashiru</string>
        <string>morning</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HASHIRU_NO_OPEN</key>
        <string>1</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/hachikenkazuki/hashiru</string>
    <key>StandardOutPath</key>
    <string>/Users/hachikenkazuki/hashiru/reports/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hachikenkazuki/hashiru/reports/launchd.log</string>
</dict>
</plist>
